# infra/resources/oom-guard.service
[Unit]
Description=SwarmFort OOM Guard
After=network.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/oom-guard.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target